# packages/digital-hydrant/digital_hydrant-0.1.7-py3-none-any.whl/digital_hydrant/systemd/digital-hydrant.service
[Unit]
Description=Digital Hydrant

[Service]
Type=simple
ExecStart=/usr/local/bin/hydrant

[Install]
WantedBy=multi-user.target
